# deploy/local-heartbeat/server-status-local-heartbeat@.service
[Unit]
Description=Server Status local heartbeat (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# A local health failure is successfully delivered as status=failed. Only
# configuration or delivery failures should leave this reporter unit failed.
SuccessExitStatus=1
User=root
Group=root
UMask=0077
ExecStart=/usr/bin/python3 /usr/local/lib/server-status-page/local_heartbeat.py --config /etc/server-status-page/local-heartbeat/%i.json
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
StandardOutput=journal
StandardError=journal
